# lapcam.service
[Unit]
Description=Laparoscopic Camera Trainer
After=multi-user.target
Wants=multi-user.target
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
User=pi
Group=pi
ExecStartPre=+/bin/sh -c 'setterm --blank force --cursor off > /dev/tty1 2>/dev/null; echo 0 > /sys/class/vtconsole/vtcon1/bind 2>/dev/null || true'
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=GPIOZERO_PIN_FACTORY=lgpio
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/pi/lapcam
ExecStart=/usr/bin/python3 /home/pi/lapcam/lapcam.py
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=42
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
